# .config/systemd/user/mnt-backups-kensington.automount
[Unit]
Description=Automount USB backup partition
ConditionPathExists=/mnt/backups/kensington

[Automount]
Where=/mnt/backups/kensington
TimeoutIdleSec=10

[Install]
WantedBy=multi-user.target
