# .config/systemd/user/remapd.service
[Unit]
Description=Run remapd to trigger remaps whenever input device is connected

[Service]
ExecStart=%h/.local/bin/remapd
